# deploy/supervisord.conf
[program:kern]
command=/opt/kern/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=/opt/kern
user=kern
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=15
redirect_stderr=true
stdout_logfile=/var/log/kern/kern.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=KERN_PRODUCT_POSTURE="production",KERN_POLICY_MODE="corporate",KERN_DB_PATH="/var/lib/kern/kern.db",KERN_SYSTEM_DB_PATH="/var/lib/kern/kern-system.db",KERN_PROFILE_ROOT="/var/lib/kern/profiles",KERN_BACKUP_ROOT="/var/lib/kern/backups",KERN_DOCUMENT_ROOT="/var/lib/kern/documents"
